# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting pgAdmin4..."

cd /pgadmin4

sudo /usr/sbin/postfix start 2>/dev/null || true

timeout=$(cd /pgadmin4 && /venv/bin/python3 -c 'import config; print(config.SESSION_EXPIRATION_TIME * 60 * 60 * 24)')

exec /venv/bin/gunicorn \
    --timeout "${timeout}" \
    --bind unix:/tmp/pgadmin4.sock \
    --workers=1 \
    --threads=25 \
    --chdir /pgadmin4 \
    -c gunicorn_config.py \
    pgAdmin4:app
