# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting pgAdmin4..."

# Run the pgAdmin4 wrapper script
exec /usr/local/bin/run-pgadmin.sh
