# pgadmin4/rootfs/etc/services.d/pgadmin4/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting pgAdmin4..."

# Run pgAdmin4 using the installed Python venv
cd /pgadmin4
exec /venv/bin/python3 /pgadmin4/run_pgadmin.py
