# dockhand_hawser/rootfs/etc/services.d/ingress/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
exec darkhttpd /var/www --port 8099 --addr 0.0.0.0 --no-listing
