# pgadmin4/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting NGINX for ingress support..."

# Start nginx in foreground
exec nginx -g "daemon off;"
