# dozzle-agent/rootfs/etc/s6-overlay/s6-rc.d/dozzle/run
#!/command/execlineb -P
foreground { echo "Starting Dozzle Docker log viewer..." }
exec /app/dozzle
